# build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D

[program:jupyter]
command=jupyter notebook --allow-root --no-browser --ip 0.0.0.0 --port 8888 --NotebookApp.token='' --NotebookApp.password=''